# abird-dn42-roa/root/etc/cron.d/cron
*/15 * * * * wget -O /bird/roa_dn42_v6.conf https://dn42.burble.com/roa/dn42_roa_bird2_6.conf
*/15 * * * * wget -O /bird/roa_dn42.conf https://dn42.burble.com/roa/dn42_roa_bird2_4.conf