# container/root/etc/services.d/babeld/run
#!/bin/sh
/usr/bin/babeld -c /abird/babeld.conf
