# abird/root/etc/cron.d/cron
*/15 * * * * /usr/sbin/birdc c
*/15 * * * * /bin/sh /ptp.sh > /etc/bird/ptp.conf