# abird/root/etc/services.d/lg/run
#!/bin/sh
/usr/local/bin/proxy --bird /var/run/bird.ctl
