# container/root/etc/cron.d/dn42roa
*/15 * * * * /roa.sh && /usr/sbin/birdc c
