# container/root/etc/services.d/cron/run
#!/bin/sh
crontab -u root /etc/cron.d/dn42roa
crond -f
