# container/root/etc/services.d/cron/run
#!/usr/bin/bash
crontab -u root /etc/cron.d/dn42roa
cron -f
