# container/root/etc/services.d/cron/run
#!/bin/sh
crontab -u root /etc/cron.d/cron
crond -f
